# etc/supervisor/conf.d/health.conf
[program:health]
command=bash -c '/dockerstartup/d.sh'
autorestart=false
