# etc/supervisor/conf.d/nginx.conf
[program:nginx]
command=bash -c '/dockerstartup/go_go.sh'
autorestart=true
directory=/root
